# units/pwupdd@.service.in
[Unit]
Description=Daemon to update passwd and shadow entries
After=local-fs.target

[Service]
ExecStart=@LIBEXECDIR@/pwupdd -d
